# dnsmasq.cron
#minute (0-59)
#|   hour (0-23)
#|   |    day of the month (1-31)
#|   |    |   month of the year (1-12 or Jan-Dec)
#|   |    |   |   day of the week (0-6 with 0=Sun or Sun-Sat)
#|   |    |   |   |   commands
#|   |    |   |   |   |
#### update dnsmasq ad server list
 15   7   *  *   1   /3 curl -s https://raw.githubusercontent.com/MLWALK3R/Automatic-dnsmasq-update/master/dnsmasq.sh  | bash>> /dev/null 2>&1
